# docs/home-dashboard-local.service
[Unit]
Description=home dashboard

[Service]
WorkingDirectory=/home/pi/home-dashboard/server
ExecStart=/usr/bin/node /home/pi/home-dashboard/server/src/Server.js
Type=simple
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
